# services/mirrors-index.path
[Unit]
Description=Mirrors index page update trigger

[Path]
PathChanged=/srv/mirror/.sync-trigger
Unit=mirrors-index.service

[Install]
WantedBy=multi-user.target
