# services/mirrors-index.timer
[Unit]
Description=Update Mirrors index page periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
